# etc/pull_battery_shutdown.crontab
* * * * *       acpi -b | grep '[^0-9][1][%]' && shutdown
