# etc/pull_battery_shutdown.crontab
* * * * *      root acpi -b | grep '[^0-9][1][\%]' && shutdown && echo "Shutting down: `acpi -b`"
